# backend/vps_supervisor_config.conf
# ================================================================
# SUPERVISOR CONFIG FOR HIGH PERFORMANCE VPS (300K+ Users)
# ================================================================
# Copy this file to: /etc/supervisor/conf.d/tg_exchange.conf
# Then run: sudo supervisorctl reread && sudo supervisorctl update
# ================================================================

[program:tg_exchange_backend]
command=/root/.venv/bin/gunicorn server:app -c gunicorn_config.py
directory=/app/backend
autostart=true
autorestart=true
user=root
environment=PYTHONUNBUFFERED="1",PYTHONDONTWRITEBYTECODE="1"
stderr_logfile=/var/log/supervisor/tg_exchange_backend.err.log
stdout_logfile=/var/log/supervisor/tg_exchange_backend.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stderr_logfile_backups=10
stdout_logfile_backups=10
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true
killasgroup=true
priority=100

[program:tg_exchange_frontend]
command=yarn start
directory=/app/frontend
autostart=true
autorestart=true
user=root
environment=HOST="0.0.0.0",PORT="3000",NODE_ENV="production"
stderr_logfile=/var/log/supervisor/tg_exchange_frontend.err.log
stdout_logfile=/var/log/supervisor/tg_exchange_frontend.out.log
stderr_logfile_maxbytes=50MB
stdout_logfile_maxbytes=50MB
stopsignal=TERM
stopwaitsecs=50
stopasgroup=true
killasgroup=true
priority=200

[group:tg_exchange]
programs=tg_exchange_backend,tg_exchange_frontend
priority=999
